# external/pulseaudio-6.0/src/daemon/systemd/user/pulseaudio.service.in
[Unit]
Description=Sound Service

[Service]
ExecStart=@PA_BINARY@ --daemonize=no
Restart=on-failure

[Install]
Also=pulseaudio.socket
WantedBy=default.target
