# external/alsa-utils-1.1.2/alsactl/90-alsa-restore.rules.in
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="controlC*", KERNELS!="card*", TEST=="@sbindir@", TEST=="@mydatadir@", GOTO="alsa_restore_go"
GOTO="alsa_restore_end"

LABEL="alsa_restore_go"
TEST!="@daemonswitch@", RUN+="@sbindir@/alsactl restore $attr{device/number}"
TEST=="@daemonswitch@", RUN+="@sbindir@/alsactl nrestore $attr{device/number}"

LABEL="alsa_restore_end"
