# rootfs/etc/inittab
# Default runlevel 5
id:5:initdefault:

# System initialization.
si1::sysinit:/sbin/ueventd.sh &
si2::sysinit:/etc/init.d/rcS

# Runlevel 5
mu:5:wait:/etc/init.d/startup.sh
con:5:respawn:/bin/getty -l /bin/login 0 console

# Runlevel 6 (reboot)
l6:6:wait:/etc/init.d/reboot.sh
